# rpmbuild/SOURCES/cron.d/disable-rdds-for-rdap-hosts-pl
0 1 1 * * root /opt/zabbix/scripts/disable-rdds-for-rdap-hosts.pl >> /var/log/zabbix/disable-rdds-for-rdap-hosts.err 2>&1
