# rpmbuild/SOURCES/cron.d/rsm-slv-rdap-downtime-pl
* * * * * root sleep 16; /opt/zabbix/scripts/slv/rsm.slv.rdap.downtime.pl >> /var/log/zabbix/rsm.slv.err 2>&1
