# rpmbuild/SOURCES/cron.d/sla-report-php
0 15 1 * * root /opt/zabbix/scripts/sla-report.php >> /var/log/zabbix/sla-report.err 2>&1
