# rpmbuild/SOURCES/cron.d/zabbix-mysql-partitioning
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

0 23 *  *  *  root    perl    /opt/zabbix/scripts/MySQL_part_management.pl >> /var/log/zabbix/zabbix-mysql-partitioning 2>&1
0 2  *  *  *  root    perl    /opt/zabbix/scripts/MySQL_part_management.pl >> /var/log/zabbix/zabbix-mysql-partitioning 2>&1
